# ykushxs_init.sh

#!/usr/bin/bash
# /etc/init.d/ykush_init.sh
### BEGIN INIT INFO
# Provides:          ykushxs_init.sh
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO
/home/USER/prisma/scripts/ykushxs_off.sh

